# matchday.service
[Unit]
Description=Matchday soccer tracker
After=network.target

[Service]
Type=simple
User=daniel
WorkingDirectory=/home/daniel/soccerscore/soccer-tracker
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
